# hvp/roles/ovirtengine/files/relay-smtps.socket
[Socket]
ListenStream=127.0.0.1:11125
Accept=yes

[Install]
WantedBy=sockets.target
