# hvp/roles/ovirtengine/files/relay-smtps@.service
[Service]
ExecStart=/usr/bin/stunnel /etc/stunnel/relay-smtps.conf
StandardInput=socket
Type=forking
User=root
PrivateTmp=true
